# Docker-Ossec/ossec-mysql/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log

[program:script]
command=service mysql restart
